# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; hermes-webui: the chat UI. Always running, exposed on $HERMES_WEBUI_PORT.
[program:webui]
command=python3 /opt/hermes-webui/server.py
directory=/opt/hermes-webui
user=hermes
autostart=true
autorestart=true
startretries=10
startsecs=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/data",HERMES_HOME="/data/.hermes",PYTHONPATH="/opt/hermes-webui"

; hermes-agent gateway: messaging platform bots (Discord, Telegram, Slack, etc.).
; Off by default. Set HERMES_GATEWAY_ENABLED=1 to activate.
; The entrypoint patches autostart=false → autostart=true via sed when enabled.
[program:gateway]
command=python3 -c "from hermes_cli.gateway import run_gateway; run_gateway()"
user=hermes
autostart=false
autorestart=true
startretries=5
startsecs=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HOME="/data",HERMES_HOME="/data/.hermes"
